# .docker/crontab
* * * * * /usr/local/bin/php /var/www/html/resources/clear-cache.php